# debian/bin/etc/cron.d/ayatana_indicator_updatenotifier
@reboot root /etc/cron.hourly/ayatana_indicator_updatenotifier
