# debian/bin/etc/cron.hourly/ayatana_indicator_updatenotifier
#!/bin/sh

MAILTO=""

test -x /usr/bin/apt-get || exit 0
/usr/bin/python3 -c 'import random; import time; time.sleep(random.random() * 180)' \
	&& /usr/bin/apt-get update
